# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/app/supervisord.pid
childlogdir=/app/logs
loglevel=info

[inet_http_server]
port = 127.0.0.1:9001

[supervisorctl]
serverurl=unix:///app/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ============================================================================
; HTTP SERVER (server.py) - Dhan Trading API wrapper
; Production: Single reliable service with auto-restart and monitoring
; ============================================================================
[program:http-server]
command=python /app/server.py
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
stopasgroup=true
priority=999
stdout_logfile=/app/logs/http-server.log
stderr_logfile=/app/logs/http-server-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=HTTP_SERVER_URL="http://localhost:3005",PYTHONUNBUFFERED="1"

; Restart policy: aggressive on startup, conservative on runtime
startsecs=10
stopwaitsecs=10

; ============================================================================
; MCP SERVER (mcp_server.py) - MCP wrapper with WebSocket support
; Runs alongside HTTP server for tool execution
; ============================================================================
[program:mcp-server]
command=python /app/mcp_server.py
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
stopasgroup=true
priority=998
stdout_logfile=/app/logs/mcp-server.log
stderr_logfile=/app/logs/mcp-server-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=HTTP_SERVER_URL="http://localhost:3005",MCP_PORT="3008",PYTHONUNBUFFERED="1"

; Restart policy: aggressive on startup, conservative on runtime
startsecs=10
stopwaitsecs=10

; ============================================================================
; Process Groups
; ============================================================================
[group:dhan-services]
programs=http-server,mcp-server
priority=999
